# terminal/systemd/claude-web-session@.service
[Unit]
Description=Persistent session daemon for claude-web t%i
After=network.target

[Service]
Type=simple
User=@APP_USER@
Group=@APP_USER@
WorkingDirectory=~
Environment=TERM=xterm-256color
Environment=LANG=en_US.UTF-8
# Point shells at the Apps xpra display so GUI (X11) apps started from a terminal
# (gnuplot, matplotlib, xeyes, ImageMagick display, …) render on the Apps desktop
# and show up as tabs in the X11 Apps app. XDG_RUNTIME_DIR keeps Qt/dbus apps
# (e.g. gnuplot's qt terminal) from warning/misbehaving.
Environment=DISPLAY=@APPS_DISPLAY@
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/@APP_UID@/bus
Environment=XDG_RUNTIME_DIR=/run/user/@APP_UID@
ExecStart=@APP_DIR@/claude-session serve %i
Restart=always
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target
