# terminal/systemd/vibetop-manager.service
[Unit]
Description=Terminal manager API (start/stop terminal instances on demand)
After=network.target
# Crash-loop guard: tolerate flapping (a bad deploy, a transient dependency) but
# give up after 10 failures in 5 min instead of hot-looping forever.
StartLimitIntervalSec=300
StartLimitBurst=10

[Service]
Type=simple
# Optional host-local overrides (absent on a default home install → no-op). The
# /opt/vibetop migration writes this to point ADMIN users + secret paths at the
# system tree: VIBETOP_ADMINS, ONLYOFFICE_SECRET_FILE, SESSION_SECRET_FILE.
EnvironmentFile=-/etc/vibetop/manager.env
ExecStart=/usr/bin/python3 @APP_DIR@/terminal-manager.py @BASE_PORT@
Restart=on-failure
RestartSec=2
# Liveness watchdog: the manager pets this only while a loopback GET /api/ping
# actually answers (see _watchdog_loop), so a WEDGED process — stuck accept
# loop, exhausted threads, deadlock — is restarted, not just a crashed one
# (which Restart=on-failure already covers). 60s window, petted every ~20s, so
# two consecutive missed probes are tolerated before a restart (no flapping
# under a brief load spike). A watchdog kill exits non-zero → Restart fires.
WatchdogSec=60

# Sandboxing — DELIBERATELY conservative. This service runs as root and is the
# one that legitimately needs broad power: it drives systemctl, drops to
# APP_USER via su/sudo (both setuid), reads sysfs/debugfs, and rewrites
# /etc/nginx + /etc/systemd + the web root during an in-app Update. So the
# usual heavy directives are OFF ON PURPOSE — NoNewPrivileges/RestrictSUIDSGID
# would break su/sudo; ProtectSystem=strict / ProtectHome would break the
# updater and home writes; PrivateTmp would hide the session sockets;
# RestrictNamespaces would break snap-confine (the manager launches snap
# Chromium for /api/browser/open and snap Firefox/Chromium for the X11 Launcher,
# and snap-confine needs to create a mount namespace — with it blocked the
# hand-off `chromium <url>` silently fails to start and the link never opens).
# What's left are the directives that add real protection without touching any
# of that: block kernel tampering, realtime abuse, and personality changes.
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectClock=yes
ProtectHostname=yes
RestrictRealtime=yes
LockPersonality=yes

[Install]
WantedBy=multi-user.target
