# files/systemd/vibetop-filebrowser.service
[Unit]
Description=FileBrowser web file manager (Vibetop)
After=network.target

[Service]
Type=simple
User=@APP_USER@
Group=@APP_USER@
WorkingDirectory=@APP_HOME@
ExecStart=@FB_BIN@ --database @FB_DB@
Restart=on-failure
RestartSec=2
# Sandboxing — match the sibling units (session, browser-xpra). FileBrowser
# never tweaks kernel knobs, loads modules, or manages cgroups, and as a plain
# Go binary it needs no privilege escalation.
NoNewPrivileges=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target
