# browser/systemd/claude-browser-xpra.service
[Unit]
Description=claude-browser xpra session (display :@DISPLAY_NUM@, HTML5 on 127.0.0.1:@XPRA_PORT@)
After=network.target

[Service]
Type=simple
User=@APP_USER@
Group=@APP_USER@
WorkingDirectory=~
Environment=HOME=@APP_HOME@
Environment=XDG_RUNTIME_DIR=/run/user/@APP_UID@
# Evict dead/stale clients after 45s instead of the 60s default, so a freshly
# opened browser doesn't wait on a zombie session from a closed tab/laptop.
# Don't go lower: phones on flaky/power-saving WiFi can stall past 20s while
# alive, and a backgrounded Safari tab stops answering pings — both got
# legitimately-connected clients evicted at the old 20s value.
Environment=XPRA_PING_TIMEOUT=45
ExecStart=/usr/bin/xpra start-desktop :@DISPLAY_NUM@ \
    --bind-ws=127.0.0.1:@XPRA_PORT@ \
    --ws-auth=none \
    --html=on \
    --daemon=no \
    --xvfb=Xorg \
    --start="matchbox-window-manager -use_titlebar no -use_desktop no" \
    --start=@LOOP_SCRIPT@ \
    --clipboard=yes \
    --clipboard-direction=both \
    --sharing=yes \
    --resize-display=yes \
    --input-devices=uinput \
    --notifications=no \
    --mdns=no \
    --printing=no \
    --pulseaudio=no \
    --speaker=off \
    --microphone=off \
    --dbus=no \
    --encoding=auto \
    --min-quality=30 \
    --quality=80 \
    --min-speed=50 \
    --speed=100 \
    --compression_level=3 \
    --video-scaling=auto \
    --sync-xvfb=0
ExecStop=/usr/bin/xpra stop :@DISPLAY_NUM@
Restart=on-failure
RestartSec=3
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target
