# claude-usage/systemd/vibetop-claude-proxy.service
[Unit]
Description=vibetop Claude usage proxy (captures Max-plan rate-limit headers)
Documentation=https://github.com/nicejunjie/vibetop
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Runs as the OPERATOR (the human whose Claude Code is observed) — NOT APP_USER,
# which may be the `vibetop` service account that never runs Claude Code. systemd
# sets $HOME to the operator's home, where the usage JSON is written and where the
# manager reads it (CLAUDE_USAGE_FILE). Loopback-only, in-place from the checkout.
User=@OPERATOR@
Environment=CLAUDE_PROXY_PORT=7690
ExecStart=@APP_DIR@/claude-usage/vibetop-claude-proxy
Restart=always
RestartSec=2
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
